# templates/borg.service.j2
# {{ansible_managed}}

[Unit]
Description=Borg Backup
 
[Service]
Type=simple
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/usr/local/bin/borg.sh
User=root
